# roles/ansible_pull/files/ansible_pull.timer
[Unit]
Description=Run ansible-pull hourly

[Timer]
OnBootSec=15min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
